# bridges/telegram_super_tma/supertma.service
[Unit]
Description=Super TMA Telegram-Matrix Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/supertma
ExecStart=/usr/bin/python3 /opt/supertma/bot.py
Restart=always
User=supertma

[Install]
WantedBy=multi-user.target 